# scripts/crontab
SHELL=/bin/bash
BASH_ENV=/container.env
*/5 * * * * /run_python.sh >> out.txt  2>&1
