# scripts/crontab
SHELL=/bin/bash
BASH_ENV=/container.env
*/5 * * * * /run_python.sh >> /dev/console  2>&1
